# com.catalyst.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.catalyst.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/python@3.12/bin/python3.12</string>
        <string>/Users/martin/CLAUDE/Claude skills/catalyst-content-analytics/scheduler.py</string>
        <string>--once</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/martin/CLAUDE/Claude skills/catalyst-content-analytics</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/martin/CLAUDE/Claude skills/catalyst-content-analytics/outputs/scheduler.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/martin/CLAUDE/Claude skills/catalyst-content-analytics/outputs/scheduler.error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
